# services/keynav.service
[Unit]
Description=Keynav

[Service]
ExecStart=keynav
Restart=always

[Install]
WantedBy=default.target
